# config/systemd/station-blanche-web.service
[Unit]
Description=Station Blanche VT - Interface Web Locale
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
User=sb-web
Group=stationblanche
WorkingDirectory=/srv/station-blanche/app
EnvironmentFile=-/etc/station-blanche/vt.env
Environment="SB_BASE_DIR=/srv/station-blanche"
Environment="SB_EVENTS_LOG=/var/log/station-blanche/events.jsonl"
ExecStart=/srv/station-blanche/app/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=3

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/station-blanche /var/log/station-blanche
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=127.0.0.1
# Adapter à l'IP du proxy si nécessaire
IPAddressAllow=192.168.50.1

[Install]
WantedBy=multi-user.target
